# debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop service if running
        if systemctl is-active --quiet talk-spotter; then
            systemctl stop talk-spotter
        fi
        # Disable service
        if systemctl is-enabled --quiet talk-spotter 2>/dev/null; then
            systemctl disable talk-spotter
        fi
        ;;
esac

#DEBHELPER#

exit 0
